# backend/config_templates/redis.conf
[program:redis]
command=redis-server
autostart=true
autorestart=true
stdout_logfile={{ app_path }}/log/redis-std.log
stderr_logfile={{ app_path }}/log/redis-err.log